# deploy/knowcodeextra.service
[Unit]
Description=KnowCodeExtra CW Test API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/knowcodeextra
ExecStart=/opt/knowcodeextra/knowcodeextra
Restart=always
RestartSec=5
# Config file location (config.toml in WorkingDirectory by default)
Environment=CONFIG_FILE=/opt/knowcodeextra/config.toml

[Install]
WantedBy=multi-user.target
